# infra/ia2.service
[Unit]
Description=IA2 PLC runtime
After=network-online.target
Wants=network-online.target

[Service]
# Atomic deploys point `current` at the latest versioned dir; the runtime
# binary + project dir always live under there.
# NB: the install root appears in several lines below (ExecStart and
# ReadWritePaths). systemd does NOT expand variables in the ExecStart
# executable path — a `${VAR}/runtime` there is a fatal unit error — so
# the path is spelled out. Installing somewhere other than /opt/ia2
# means editing every occurrence (and the edge's configured
# `install_dir`, which deploy checks for drift).
# `--static-dir` serves the standalone HMI panel at /hmi off the deployed
# web assets (deploy places them at current/web; harmless if absent).
# The default bind is loopback, matching the ssh-tunnel trust perimeter —
# operators on the plant network need either a tunnel or an explicit
# decision to widen `--bind` (which also exposes /write and /stop; put a
# reverse proxy or firewall in front if you do).
ExecStart=/opt/ia2/current/runtime \
  --project-dir /opt/ia2/current/project \
  --bind 127.0.0.1:13001 \
  --static-dir /opt/ia2/current/web

# Restart on crash but back off if we're crash-looping (don't hammer the
# bus / journald).
Restart=on-failure
RestartSec=3s
StartLimitIntervalSec=60
StartLimitBurst=5

# Stop with SIGTERM first; the runtime stops the scan loop cooperatively.
# 10s should be plenty even for slow Modbus connect timeouts.
KillSignal=SIGTERM
TimeoutStopSec=10s

# Logs go to journald — `journalctl -u ia2 -f` on the edge.
StandardOutput=journal
StandardError=journal

# Run as root for now so EtherCAT (when wired) has CAP_NET_RAW. For Modbus-
# only deployments, switch to a dedicated user:
#   User=ia2
#   AmbientCapabilities=
# and remove the cap line below.
User=root
AmbientCapabilities=CAP_NET_RAW

# Hardening — tighten what the runtime can reach. Loosen if you find a
# legitimate access blocked.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
# Allow writes only under the install dir; everything else is read-only.
ReadWritePaths=/opt/ia2

[Install]
WantedBy=multi-user.target
